# templates/default/kafka.service.erb
[Unit]
<% @service_unit.each do |k, v| -%>
<%= k %>=<%= v %>
<% end -%>

[Service]
User=<%= @user %>
ExecStart=<%= @start_script %> <%= @configuration %>
Restart=always
RestartSec=5
